# misc/10-local.rules
KERNEL=="sd[b-z]", ACTION=="add", RUN+="/usr/bin/systemctl start cs-add-dev@%k.service"
