# cs-add-dev@.service
[Unit]
Description=CopyStation device added service

[Service]
Type=oneshot
ExecStart=/home/copycat/copystation/attach.sh %I
